# resources/configs/supervisord.conf
# PUT THIS FILE IN /etc/supervisord.conf

[supervisord]
logfile=/home/ec2-user/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10

[unix_http_server]
file=/etc/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///etc/supervisor.sock

[program:annotations]
command=/bin/bash -c 'resources/scripts/prod.sh'
directory=/home/ec2-user/annotations
autostart=true
autorestart=true
startretries=10
stderr_logfile=/home/ec2-user/logs/supervisor_annotations_error.log
stdout_logfile=/home/ec2-user/logs/supervisor_annotations.log
user=ec2-user
environment=PYTHONPATH='.:vendor',PYPY_GC_MAX='6GB',HOME='/home/ec2-user'
